# crontab/audit.cron
#!/bin/bash
export DJANGO_SETTINGS_MODULE=tuangou.settings
python2.6 /var/www/tuangou/manage.py audit_uv_pv
